# ljk_slack_bot.service
[Unit]
Description=LJK Slack Bot

[Service]
ExecStart=/usr/local/bin/run_ljk_slack_bot.sh

[Install]
WantedBy=multi-user.target
